# profile.d/history.sh
#!/bin/bash

HISTTIMEFORMAT="%m_%d %H:%M "
export HISTTIMEFORMAT
